# deploy/lernmanager.service
[Unit]
Description=Lernmanager Flask Application
After=network.target

[Service]
Type=simple
User=lernmanager
Group=lernmanager
WorkingDirectory=/opt/lernmanager

# Load environment variables (including secrets) from file
# This file is created during setup and persists across updates
EnvironmentFile=/opt/lernmanager/.env

# Temporary directory for file uploads (Flask/Werkzeug needs this)
Environment="TMPDIR=/opt/lernmanager/instance/tmp"
# Python bytecode cache directory (improves startup time)
Environment="PYTHONPYCACHEPREFIX=/opt/lernmanager/instance/tmp"

ExecStart=/opt/lernmanager/venv/bin/python run.py
Restart=always
RestartSec=5

# Security hardening (optional but recommended)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/lernmanager/data /opt/lernmanager/instance/uploads /opt/lernmanager/instance/tmp

[Install]
WantedBy=multi-user.target
