# deploy/lernmanager.service
[Unit]
Description=Lernmanager Flask Application
After=network.target

[Service]
Type=simple
User=lernmanager
Group=lernmanager
WorkingDirectory=/opt/lernmanager

# Generate secure keys: python3 -c "import secrets; print(secrets.token_hex(32))"
Environment="SECRET_KEY=CHANGE_ME_TO_RANDOM_STRING"
Environment="FLASK_ENV=production"
# Temporary directory for file uploads (Flask/Werkzeug needs this)
Environment="TMPDIR=/opt/lernmanager/instance/tmp"
# Python bytecode cache directory (improves startup time)
Environment="PYTHONPYCACHEPREFIX=/opt/lernmanager/instance/tmp"
# Optional: Database encryption key (requires sqlcipher3-binary)
# Environment="SQLCIPHER_KEY=CHANGE_ME_TO_RANDOM_STRING"

ExecStart=/opt/lernmanager/venv/bin/python run.py
Restart=always
RestartSec=5

# Security hardening (optional but recommended)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/lernmanager/data /opt/lernmanager/instance/uploads /opt/lernmanager/instance/tmp

[Install]
WantedBy=multi-user.target
