# deploy/lernmanager.service
[Unit]
Description=Lernmanager Flask Application
After=network.target

[Service]
Type=simple
User=lernmanager
Group=lernmanager
WorkingDirectory=/opt/lernmanager

# Generate a secure key: python3 -c "import secrets; print(secrets.token_hex(32))"
Environment="SECRET_KEY=CHANGE_ME_TO_RANDOM_STRING"

ExecStart=/opt/lernmanager/venv/bin/python run.py
Restart=always
RestartSec=5

# Security hardening (optional but recommended)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/lernmanager/data /opt/lernmanager/static/uploads

[Install]
WantedBy=multi-user.target
